# deploy/files/lifegame.supervisord.conf
[program:lifegame]
command= /opt/lifegame/main
directory = /opt/lifegame
autorestart=true
redirect_stderr=true
stdout_logfile = /var/log/lifegame/lifegame.log
